# env-sensors-scanner.service
[Unit]
Description=BTLE environmental sensor advertisment scanner
After=network.target

[Service]
ExecStart=/home/pi/env-stats/env/bin/python3 /home/pi/env-stats/scan_thermometer.py
WorkingDirectory=/home/pi/env-stats
StandardOutput=inherit
StandardError=inherit
Restart=always
User=pi

[Install]
WantedBy=multi-user.target
